# scripts/openrc-m5755-server
#!/sbin/openrc-run
# 5755 平台服务端 OpenRC 服务(Alpine,CTID 105/106 实际运行形态)。
# 安装:复制为 /etc/init.d/m5755-server,chmod +x,rc-update add m5755-server default。
# 二进制 /opt/m5755/m5755-server(deploy.sh 原子替换),环境文件 /opt/m5755/.env。
# 注:仓库内 m5755-server.service 是 systemd 单元,仅供 systemd 宿主;OpenRC 宿主用本文件。

name="m5755-server"
description="5755 platform server (SDK gateway v2)"
command="/opt/m5755/m5755-server"
command_user="m5755:m5755"
command_background=true
pidfile="/run/m5755-server.pid"
output_log="/var/log/m5755-server.log"
error_log="/var/log/m5755-server.log"

depend() {
    need net
    after firewall
}

start_pre() {
    if [ ! -f /opt/m5755/.env ]; then
        eerror "缺少 /opt/m5755/.env(DATABASE_URL/密钥/PLATFORM_ENV)"
        return 1
    fi
    # 把环境文件导入服务启动环境(production 缺密钥会 fail-closed 拒绝启动)
    set -a
    . /opt/m5755/.env
    set +a
}
